# scripts/tokymon.service
[Unit]
Description=Tokymon Robot Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/ankursharma/Projects/tokymon
ExecStart=/home/ankursharma/Projects/tokymon/scripts/run_tokymon.sh
EnvironmentFile=/home/ankursharma/Projects/tokymon/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
